# fevr/rootfs/etc/services.d/fevr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
#declare message
declare frigateUrl
declare sqlitewebPort
## Get the 'message' key from the user config options.
frigate=$(bashio::config 'frigate')
sqlitewebPort=$(bashio::config 'sqlitewebPort')

## Run your program
exec sed -i '/http:\/\/frigate.mydomain.com/c\http://${frigateUrl}' /var/www/config/config.json && \
     sqlite_web -p ${sqlitewebPort} -H 0.0.0.0 -x /var/www/db/fEVR.sqlite